# trunk/install/runworker.conf
[program:djangoworkders]
directory = /var/www/html/onlineTest
command = python3 ./manage.py runworker --threads 2
autostart = true
autorestart = true
user = judge
redirect_stderr = true

